# Automation/packer/exporter.ini
[program:wowza_exporter]
command=python38 customExporter.py
directory=/usr/src/wolbf/Exporter/app
stdout_logfile=/tmp/test_process_output.txt
redirect_stderr=true

